# deploy/staging/metaforge-api-staging.service
[Unit]
Description=Metaforge Go API (staging)
After=network.target

[Service]
Type=simple
User=agent
WorkingDirectory=__WORKTREE__
ExecStart=__WORKTREE__/api/metaforge \
    -port 8081 \
    -db __WORKTREE__/data-pipeline/output/lexicon_v2.db \
    -strings __WORKTREE__/strings \
    -cors-origin https://metaforge-next.julianit.me
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=metaforge-api-staging
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=__WORKTREE__/data-pipeline/output

[Install]
WantedBy=multi-user.target
